# luci-app-memos/root/etc/init.d/memos
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-2.0-only
#
# memos 服务的 procd init 脚本。
# 启动时若 binpath 不存在，后台运行 memos.sh download，完成后由 init 再 start。

USE_PROCD=1
START=95
STOP=01

CONFIGURATION=memos
UPDATE_SCRIPT=/usr/share/memos/memos.sh
CRON_FILE=/etc/crontabs/root
INSTANCE_NAME=memos

extra_command "isrunning" "report whether memos is running"

memos_get() { uci -q get "${CONFIGURATION}.${CONFIGURATION}.$1" 2>/dev/null; }
memos_getd() { local v; v="$(memos_get "$1")"; echo "${v:-$2}"; }

# 每行一个参数，避免路径/DSN 含空格时被拆开
build_args_file() {
	local out="$1"
	local addr port data driver dsn instance_url unix_sock log_level demo allow_private

	: > "$out"

	addr=$(memos_get addr)
	port=$(memos_getd port 8081)
	data=$(memos_getd data /etc/memos)
	driver=$(memos_getd driver sqlite)
	dsn=$(memos_get dsn)
	instance_url=$(memos_get instance_url)
	unix_sock=$(memos_get unix_sock)
	log_level=$(memos_getd log_level info)
	demo=$(memos_get demo)
	allow_private=$(memos_get allow_private_webhooks)

	if [ -n "$unix_sock" ]; then
		printf '%s\n' --unix-sock "$unix_sock" >> "$out"
	else
		[ -n "$addr" ] && printf '%s\n' --addr "$addr" >> "$out"
		printf '%s\n' --port "$port" >> "$out"
	fi

	printf '%s\n' --data "$data" >> "$out"
	printf '%s\n' --driver "$driver" >> "$out"

	if [ "$driver" != "sqlite" ] && [ -n "$dsn" ]; then
		printf '%s\n' --dsn "$dsn" >> "$out"
	fi

	[ -n "$instance_url" ] && printf '%s\n' --instance-url "$instance_url" >> "$out"
	[ "$demo" = "1" ] && printf '%s\n' --demo >> "$out"
	[ "$allow_private" = "1" ] && printf '%s\n' --allow-private-webhooks >> "$out"
	printf '%s\n' --log-level "$log_level" >> "$out"
}

sync_cron() {
	local default_line enabled autoupdate tmp
	enabled=$(memos_getd enabled 0)
	autoupdate=$(memos_getd cron_autoupdate 0)
	default_line="30 3 * * * $UPDATE_SCRIPT update >/dev/null 2>&1"

	mkdir -p "$(dirname "$CRON_FILE")" 2>/dev/null
	[ -f "$CRON_FILE" ] || : > "$CRON_FILE"

	tmp="${CRON_FILE}.memos.$$"
	grep -vF "$UPDATE_SCRIPT" "$CRON_FILE" > "$tmp" 2>/dev/null || : > "$tmp"
	mv -f "$tmp" "$CRON_FILE"

	if [ "$enabled" = "1" ] && [ "$autoupdate" = "1" ]; then
		echo "$default_line" >> "$CRON_FILE"
	fi
	log_cron_restart
}

log_cron_restart() {
	/etc/init.d/cron enabled 2>/dev/null && /etc/init.d/cron restart >/dev/null 2>&1 || true
}

start_service() {
	local enabled binpath data

	enabled=$(memos_getd enabled 0)
	binpath=$(memos_getd binpath /usr/bin/memos)
	data=$(memos_getd data /etc/memos)

	sync_cron

	[ "$enabled" = "0" ] && return 0

	mkdir -p "$data" 2>/dev/null

	if [ ! -x "$binpath" ]; then
		logger -t memos "未找到二进制 $binpath，后台下载中"
		(
			"$UPDATE_SCRIPT" download >/dev/null 2>&1 && {
				logger -t memos "下载完成，启动服务"
				/etc/init.d/memos start
			} || logger -t memos "后台下载失败，见 /tmp/memos_update.log"
		) &
		return 0
	fi

	local user group gc maxprocs memlimit tz args_file arg
	user=$(memos_get user)
	group=$(memos_get group)
	gc=$(memos_getd gc 0)
	maxprocs=$(memos_getd maxprocs 0)
	memlimit=$(memos_getd memlimit 0)

	tz=""
	[ -f /etc/TZ ] && tz=$(cat /etc/TZ 2>/dev/null)
	[ -z "$tz" ] && [ -f /var/etc/TZ ] && tz=$(cat /var/etc/TZ 2>/dev/null)

	args_file="/tmp/memos-procd-args.$$"
	build_args_file "$args_file"

	procd_open_instance "$INSTANCE_NAME"
	procd_set_param command "$binpath"
	while IFS= read -r arg || [ -n "$arg" ]; do
		[ -n "$arg" ] && procd_append_param command "$arg"
	done < "$args_file"
	rm -f "$args_file"

	procd_set_param respawn 3600 5 5
	procd_set_param stderr 1
	procd_set_param limits core=unlimited nofile=65535
	[ -n "$user" ]  && procd_set_param user  "$user"
	[ -n "$group" ] && procd_set_param group "$group"
	[ -n "$tz" ] && procd_set_param env TZ="$tz"
	[ "$gc" != "0" ]       && procd_set_param env GOGC="$gc"
	[ "$maxprocs" != "0" ] && procd_set_param env GOMAXPROCS="$maxprocs"
	[ "$memlimit" != "0" ] && procd_set_param env GOMEMLIMIT="${memlimit}MiB"
	procd_set_param file /etc/config/memos
	procd_close_instance

	logger -t memos "服务已启动 (bin=$binpath)"
}

reload_service() { start; }

stop_service() { sync_cron; logger -t memos "服务已停止"; }

service_triggers() {
	procd_add_reload_trigger "$CONFIGURATION"
}

isrunning() {
	local binpath
	binpath=$(memos_getd binpath /usr/bin/memos)
	[ -x "$binpath" ] || { echo "no bin"; return 2; }
	if pgrep -x "${binpath##*/}" >/dev/null 2>&1 || pgrep -f "$binpath" >/dev/null 2>&1; then
		echo "running"
		return 0
	fi
	echo "not run"
	return 1
}

boot() {
	rc_procd start
}
