# luci-app-memos/root/etc/uci-defaults/99-luci-app-memos
#!/bin/sh

# 创建默认数据目录，供 memos 首次写入 SQLite
mkdir -p /etc/memos 2>/dev/null

# 注册 ucitrack，配置变更时重载服务
uci -q get ucitrack.@memos[-1] >/dev/null 2>&1 || {
	uci -q add ucitrack memos >/dev/null 2>&1
	uci -q set ucitrack.@memos[-1].init=memos >/dev/null 2>&1
	uci -q commit ucitrack >/dev/null 2>&1
}

exit 0
